# provision/app/etc/supervisord.d/mailcatcher.ini
[program:mailcatcher]
command=/bin/bash -c "mailcatcher --smtp-port 25 --ip $(ip addr show enp0s8 | grep inet | head -1 | cut -d' ' -f6 | cut -d'/' -f1)"
